# templates/ubuntu/nginx.service.erb
[Unit]
Description=A high performance web server and a reverse proxy server
After=network.target

[Service]
Type=forking
PIDFile=<%= node['nginx']['pid'] %>
ExecStartPre=<%= node['nginx']['binary'] %> -t -g 'daemon on; master_process on;'
ExecStart=<%= node['nginx']['binary'] %> -g 'daemon on; master_process on;'
ExecReload=<%= node['nginx']['binary'] %> -g 'daemon on; master_process on;' -s reload
ExecStop=-/sbin/start-stop-daemon --quiet --stop --retry QUIT/5 --pidfile <%= node['nginx']['pid'] %>
TimeoutStopSec=5
KillMode=mixed

[Install]
WantedBy=multi-user.target
